# trading-bot.service
# P2修复: systemd 服务配置文件
# 用于 Linux VPS 24h 守护运行
# 
# 安装步骤:
# 1. 复制到 /etc/systemd/system/trading-bot.service
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable trading-bot
# 4. sudo systemctl start trading-bot
#
# 查看日志: journalctl -u trading-bot -f

[Unit]
Description=Quantitative Trading Bot Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-bot
ExecStart=/home/ubuntu/trading-bot/.venv/bin/python separated_system/trade_engine.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/trading-bot/logs/systemd.log
StandardError=append:/home/ubuntu/trading-bot/logs/systemd-error.log

# 环境变量
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/ubuntu/trading-bot/.env

# 资源限制
MemoryMax=1G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
